# zfs/etc/systemd/system/stmf-svc.service.in
[Unit]
Description=Start stmf service
DefaultDependencies=no
After=zfs.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStartPre=sleep 2
ExecStart=@sbindir@/stmfsvc start

[Install]
WantedBy=multi-user.target
