# zfs/etc/systemd/system/zfs-clusterd.service.in
[Unit]
Description=Start clusterd service
DefaultDependencies=no
Requires=NetworkManager.service
After=zfs.target

[Service]
ExecStart=@sbindir@/clusterd -d
Restart=on-failure

[Install]
WantedBy=multi-user.target
