# zfs/etc/systemd/system/zfs-boot.service.in
[Unit]
Description=Probe zfs module
DefaultDependencies=no
Requires=systemd-udev-settle.service
After=systemd-udev-settle.service
After=cryptsetup.target
After=systemd-remount-fs.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/sbin/modprobe zfs

[Install]
WantedBy=zfs.target
